# recipes-devtools/rknpu2/rockchip-npu2-runtime/rknn-server.service.in
[Unit]
Description=Rockchip NPU service

[Service]
Type=simple
ExecStart=@BINDIR@/rknn_server
Restart=always

[Install]
WantedBy=multi-user.target
