# recipes-bsp/rkwifibt/rkwifibt-idea3588/rk-bluetooth-init.service.in
[Unit]
Description=Bluetooth module initialization
Before=bluetooth.service

[Service]
Type=simple
ExecStart=@SBINDIR@/rk-bluetooth-init

[Install]
WantedBy=bluetooth.target
